# ntp.service
[Unit]
After=network.target

[Service]
ExecStart=/bin/sleep inf
Restart=always
Type=simple
